# launchd/com.ssoward.slack-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ssoward.slack-bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/ssoward/sandbox/workspace/slack-bot/dist/index.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ssoward/sandbox/workspace/slack-bot</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/ssoward/sandbox/workspace/slack-bot/slack-bot.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ssoward/sandbox/workspace/slack-bot/slack-bot.err.log</string>
</dict>
</plist>
